# meta-demo/recipes-core/data-layout/files/data-layout.service
[Unit]
Description=Prepare persistent demo data directories
After=local-fs.target
RequiresMountsFor=/data

[Service]
Type=oneshot
ExecStart=/usr/libexec/demo/data-layout.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
