# meta-demo/recipes-security/firewall/files/firewall.service
[Unit]
Description=Demo firewall baseline
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/demo-firewall
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
